# deploy/ledger-chain-watcher.service
[Unit]
Description=Telegram Ledger Chain Watcher
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# 高级调度参数与密码、API Key 分开保存。
EnvironmentFile=-/etc/ledger-chain-watcher/config/defaults.env
EnvironmentFile=/etc/ledger-chain-watcher/env
ExecStart=/usr/local/bin/ledger-chain-watcher
Restart=always
RestartSec=5

# 如需使用非 root 用户运行，请先创建用户并确认该用户能读取：
# /etc/ledger-chain-watcher/env
# /etc/ledger-chain-watcher/config/defaults.env
# 然后取消下面两行注释并替换成真实用户/用户组。
# User=ledger
# Group=ledger

NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
